# share/systemd/ams-isoftstat.service
[Unit]
Description=AMS iSoftrepo Service

[Service]
Type=simple
ExecStart=/opt/ams-servers/8.0/sbin/ams-isoftrepo

[Install]
WantedBy=multi-user.target
